# deploy/astroai-api.service
[Unit]
Description=AstroAI API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/astroai/server
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
